# gowebserver_1.0/DEBIAN/postinst
#! /bin/bash
update-rc.d gowebserver defaults
if which invoke-rc.d >/dev/null 2>&1; then
	invoke-rc.d gowebserver start
else
	/etc/init.d/gowebserver start
fi
